# host-freebsd/tarball/etc/periodic/daily/403.ae3
#!/bin/sh

# update NTP date NOW!
ntpdate -u -b pool.ntp.org || true

AE3_DIR="/usr/local}/share/ae3"

exec "$(which bash)" -e "$AE3_DIR/bin/periodic/daily"
