# host-freebsd/tarball/etc/periodic/daily/403.ae3
#!/bin/sh

# update NTP date NOW!
ntpdate -u -b pool.ntp.org || true

AE3_DIR="$(echo ${PKG_PREFIX-/usr/local}/share/ae3 | /usr/bin/sed -e 's|//|/|g')"

exec "$(which bash)" -e "$AE3_DIR/bin/periodic/daily"
