# host-freebsd/tarball/etc/periodic/daily/403.ae3
#!/bin/sh

# update NTP date NOW!
ntpdate -u -b pool.ntp.org || true

AE3_DIR="/usr/local/share/ae3"

export PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin:
exec bash -e "$AE3_DIR/bin/periodic/daily"
